# config/trend.cron
*/1 * * * * /home/pi/MinerApp/pi-miner/api/analytics.sh 30 2