# config/trend.cron
*/1 * * * * /Users/johngrese/Projects/MinerApp/pi-miner/api/analytics.sh 30 2